# indigo/etc/catkin/profile.d/15.rocon_launch.zsh
# generated from rocon_launch/env-hooks/15.rocon_launch.zsh.em

. "/opt/ros/indigo/share/rocon_launch/shells/env.zsh"
