# indigo/etc/catkin/profile.d/15.rocon_launch.tcsh
# generated from rocon_launch/env-hooks/15.rocon_launch.tcsh.em

. "/opt/ros/indigo/share/rocon_launch/shells/env.tcsh"
